# deploy/ddhotel-bot.service

[Unit]
Description=DDHotel Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bots/ddhotel-test
Environment=PATH=/root/bots/ddhotel-test/.venv/bin:/usr/bin:/bin
ExecStart=/root/bots/ddhotel-test/.venv/bin/python main.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
