# deploy/ddhotel-bot.service
# Установка (VPS):
#   sudo cp deploy/ddhotel-bot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ddhotel-bot
# Лог: journalctl -u ddhotel-bot -f
# Другой пользователь/путь — поправьте User=, WorkingDirectory и ExecStart.

[Unit]
Description=DDHotel Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bots/ddhotel-test
Environment=PATH=/root/bots/ddhotel-test/.venv/bin:/usr/bin:/bin
ExecStart=/root/bots/ddhotel-test/.venv/bin/python main.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
